# rituals/launchd/com.alpha.assistant.morning.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.alpha.assistant.morning</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>open -a "Claude" --args --message "/ritual morning — run the morning check-in for <AI_OS_FOLDER>"</string>
    </array>

    <!-- Fire every weekday at HH:MM (user's local time) -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Weekday</key><integer>1</integer>
            <key>Hour</key><integer><HH></integer>
            <key>Minute</key><integer><MM></integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>2</integer>
            <key>Hour</key><integer><HH></integer>
            <key>Minute</key><integer><MM></integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>3</integer>
            <key>Hour</key><integer><HH></integer>
            <key>Minute</key><integer><MM></integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>4</integer>
            <key>Hour</key><integer><HH></integer>
            <key>Minute</key><integer><MM></integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>5</integer>
            <key>Hour</key><integer><HH></integer>
            <key>Minute</key><integer><MM></integer>
        </dict>
    </array>

    <!-- Log stdout/stderr to the folder so the assistant can inspect failures -->
    <key>StandardOutPath</key>
    <string><AI_OS_FOLDER>/logs/ritual-morning.out.log</string>
    <key>StandardErrorPath</key>
    <string><AI_OS_FOLDER>/logs/ritual-morning.err.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
